# virtualhere-supervisor/rootfs/etc/s6-overlay/s6-rc.d/virtualhere-supervisor/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare hostname
declare usb_device_id

#Load Config
hostname=$(bashio::config 'hostname')
usb_device_id=$(bashio::config 'usb_device_id')


#Initalize USBIP Kernel Module
modprobe vhci-hcd

## Run your program
docker run -t --privileged --rm --name vhworker --env VH_HOSTNAME="${hostname}" --env VH_USB_DEVICE_ID="${usb_device_id}" wantmys2000/virtualhere-ha-worker:latest
