# virtualhere-supervisor/rootfs/etc/s6-overlay/s6-rc.d/virtualhere-supervisor/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare hostname
declare usb_device_id

hostname=$(bashio::config 'hostname')
usb_device_id=$(bashio::config 'usb_device_id')

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "${hostname}"

## Run your program
echo "Hello there"
echo "${usb_device_id}"
echo "Device"

docker run -t --privileged --rm --name vhworker --env VH_HOSTNAME="${hostname}" --env VH_USB_DEVICE_ID="${usb_device_id}" wantmys2000/virtualhere-ha-worker:latest
